# host/agy-proxy.service
[Unit]
Description=AGY CLI Proxy - Bridges Docker to host agy CLI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root

WorkingDirectory=/opt/bb-private-kitchen/host
ExecStart=/root/.local/bin/uv run python3 agy_proxy.py

Restart=always
RestartSec=5

Environment=PYTHONUTF8=1
Environment=AGY_PROXY_HOST=0.0.0.0
Environment=AGY_PROXY_PORT=8765
Environment=AGY_PROXY_TIMEOUT=120

[Install]
WantedBy=multi-user.target
